# autoinstall/custom/systemd/ptp4l.service.d/override.conf
[Unit]
Description=Start ptp4l service
Requires=network.target
After=network.target

[Service]
ExeStart=/usr/sbin/ptp4l -f /etc/linuxptp/ptp4l.conf -i eth0

[Install]
WantedBy=multi-user.target
